# supervisord.conf
[supervisord]
nodaemon=true
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
[supervisorctl]

[program:node]
logfile = /tmp/node.log
command=yarn dev
